# get-of-metrics/DEBIAN/postinst
systemctl enable get-of-metrics
systemctl start get-of-metrics
systemctl is-active --quiet get-of-metrics && echo "Service is enabled, and running. You can see the logs with [journalctl -b -u get-of-metrics]"