# get_of_metrics/DEBIAN/prerm
sudo systemctl stop get_of_metrics
sudo systemctl disable get_of_metrics
sudo rm /etc/systemd/system/get_of_metrics.service
echo "Service is disabled and service file has been removed"